# backupd.conf
# backupd -- backup to external disk when plugged in

description	"automatically back up to external disk"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec backupd
